# services/stay-compass-kiosk.service
[Unit]
Description=Stay Compass Kiosk Compatibility Stub
After=getty@tty1.service
Requires=getty@tty1.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
