# services/stay-compass-splash.service
[Unit]
Description=Stay Compass Boot Splash Handoff
DefaultDependencies=no
After=local-fs.target plymouth-start.service
Before=stay-compass-kiosk.service
Wants=plymouth-start.service

[Service]
Type=oneshot
ExecStart=/usr/bin/plymouth show-splash
ExecStart=/usr/bin/plymouth change-mode --boot-up
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
